# bentopdf/ops/bentopdf.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=bentopdf
Image=ghcr.io/ngarside/bentopdf:1.7.4.41
Memory=8g
Network=bentopdf.network:ip=10.23.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
