# bentopdf/ops/bentopdf.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=bentopdf
DropCapability=all
Image=ghcr.io/ngarside/bentopdf:2.8.4.517@sha256:deceda9904bef58df16b76feb3d15f3d6dfe7580712efe5be745d25e60bc8e8d
Memory=512m
Network=bentopdf.network:ip=10.23.0.2
NoNewPrivileges=true
PidsLimit=64
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
