# bentopdf/ops/bentopdf.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=bentopdf
DropCapability=all
Image=ghcr.io/ngarside/bentopdf:2.8.2.471
Memory=8g
Network=bentopdf.network:ip=10.23.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
