# bentopdf/ops/bentopdf.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=bentopdf
DropCapability=all
Image=ghcr.io/ngarside/bentopdf:2.8.5.578@sha256:7a3353cf27ae72ce7405b97b56a5dd797cd4828b1960e11a80837894f459fcae
Memory=512m
Network=bentopdf.network:ip=10.23.0.2
NoNewPrivileges=true
PidsLimit=64
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
