# excalidraw/ops/excalidraw.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=excalidraw
Image=ghcr.io/ngarside/excalidraw:0.0.0.94
Memory=8g
Network=excalidraw.network:ip=10.18.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
