# excalidraw/ops/excalidraw.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=excalidraw
DropCapability=all
Image=ghcr.io/ngarside/excalidraw:0.0.0.614@sha256:e98ad753a2f3daacb61d32b4134cf17b1c36b6a901cb776bf1df254dc7153979
Memory=512m
Network=excalidraw.network:ip=10.18.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
