# excalidraw/ops/excalidraw.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=excalidraw
DropCapability=all
Image=ghcr.io/ngarside/excalidraw:0.0.0.445@sha256:67dc2c31cd276870ef24e3e0a7a11851ccb2db5153ca674130cab6a39e5137c4
Memory=512m
Network=excalidraw.network:ip=10.18.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
