# whoami/ops/whoami.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=whoami
Image=ghcr.io/ngarside/whoami:1.11.0.7
Memory=8g
Network=whoami.network:ip=10.20.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
