# fossflow/ops/fossflow.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=fossflow
DropCapability=all
Image=ghcr.io/ngarside/fossflow:0.0.0.650@sha256:e5705afa6a717e08146e067cbb5fc2aeaf1db9752150d6dc5b0845560eab343b
Memory=512m
Network=fossflow.network:ip=10.1.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
