# fossflow/ops/fossflow.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=fossflow
DropCapability=all
Image=ghcr.io/ngarside/fossflow:0.0.0.580@sha256:768625c5b45cf2490d8fe8ae1bafc046c86455cb32f43c8e5b50ca8b25e9153c
Memory=512m
Network=fossflow.network:ip=10.1.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
