# fossflow/ops/fossflow.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=fossflow
Image=ghcr.io/ngarside/fossflow:0.0.0.194
Memory=8g
Network=fossflow.network:ip=10.1.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
