# fossflow/ops/fossflow.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=fossflow
DropCapability=all
Image=ghcr.io/ngarside/fossflow:0.0.0.716@sha256:ad33f323483bbe3f4cfa18f9cef16c0c618d0060172f9befeb755fac36f2da0e
Memory=512m
Network=fossflow.network:ip=10.1.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
