# caddy/ops/caddy.container
# This is free and unencumbered software released into the public domain.

[Unit]

[Container]
ContainerName=caddy
Image=ghcr.io/ngarside/caddy:latest
Network=adguardhome.network:ip=10.2.0.3
Network=fossflow.network:ip=10.1.0.3
Network=penpot_frontend.network:ip=10.10.0.2
PublishPort=80:80/tcp
ReadOnly=true
UserNS=auto:size=8
Volume=/etc/caddy:/etc/caddy:ro

[Service]
Restart=always
TimeoutStartSec=300

[Install]
WantedBy=default.target
