# chartdb/ops/chartdb.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=chartdb
DropCapability=all
Image=ghcr.io/ngarside/chartdb:1.20.1.577@sha256:cf08fe649af617ca64dcde1942d273f53982ede95fe46ad7819e37d223f87935
Memory=512m
Network=chartdb.network:ip=10.16.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
