# chartdb/ops/chartdb.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=chartdb
Image=ghcr.io/ngarside/chartdb:1.17.0.34
Memory=8g
Network=chartdb.network:ip=10.16.0.2
NoNewPrivileges=true
PodmanArgs=--cpus=1
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
