# chartdb/ops/chartdb.container
# This is free and unencumbered software released into the public domain.

[Unit]

[Container]
ContainerName=chartdb
Image=ghcr.io/ngarside/chartdb:1.16.0.2
Network=chartdb.network:ip=10.16.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
