# chartdb/ops/chartdb.container
# This is free and unencumbered software released into the public domain.

[Container]
AddCapability=CAP_NET_BIND_SERVICE
ContainerName=chartdb
DropCapability=all
Image=ghcr.io/ngarside/chartdb:1.20.1.419@sha256:944c7fb04ae8003241eb22d0771bf5a31b9685420b47d1abf9a171ba0a7c4892
Memory=512m
Network=chartdb.network:ip=10.16.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
CPUQuota=100%
Restart=always
TimeoutStartSec=900

[Install]
WantedBy=default.target
