# vaultwarden/ops/vaultwarden.container
# This is free and unencumbered software released into the public domain.

[Unit]

[Container]
ContainerName=vaultwarden
Image=docker.io/vaultwarden/server:1.34.3-alpine
Network=vaultwarden.network:ip=10.12.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8
Volume=/var/data/vaultwarden/data:/data:U,Z

[Service]
Restart=always
TimeoutStartSec=300
ExecStartPre=mkdir --parents /var/data/vaultwarden/data

[Install]
WantedBy=default.target
