# penpot/ops/valkey.container
# This is free and unencumbered software released into the public domain.

[Container]
ContainerName=penpot_valkey
Image=ghcr.io/ngarside/valkey:8.1.4.143
Network=penpot_backend.network:ip=10.4.0.2
NoNewPrivileges=true
ReadOnly=true
UserNS=auto:size=8

[Service]
Restart=always
TimeoutStartSec=9000

[Install]
WantedBy=default.target
